# contrib/mqtt-wss/debian/mqtt-wss/DEBIAN/postinst
#!/bin/sh

update-rc.d mqtt-wss defaults
invoke-rc.d mqtt-wss start
