# system_rules/debian/postinst
#!/bin/sh
set -e

invoke-rc.d wb-rules restart

#DEBHELPER#

exit 0
